# nspawn/sample.service
[Service]
ExecStartPost=/usr/local/etc/nspawn/machines/sample.sh
