# nspawn/usr/local/etc/nspawn/profile.d/xorg.sh
#!/bin/sh

export DISPLAY=:0
export XAUTHORITY=/tmp/nspawn/xauth
export PULSE_SERVER="0.0.0.0:4656"
